# projects/Amlogic/filesystem/usr/lib/systemd/system/khadas-fan.service
[Unit]
Description=Khadas-fan

[Service]
Type=simple
ExecStart=/opt/fan

[Install]
WantedBy=basic.target
